# agents/rsrch/deploy/rsrch-server@.service
[Unit]
Description=Rsrch Research Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Obsi/Prods/01-pwf/agents/rsrch
ExecStart=/usr/bin/node /home/%i/Obsi/Prods/01-pwf/agents/rsrch/dist/server.js
Restart=on-failure
RestartSec=10
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
